# deploy_tools/gunicorn-systemd.template.socket
[Unit]
Description=gunicorn-DOMAIN-staging socket

[Socket]
ListenStream=/run/gunicorn-DOMAIN-staging.socket

[Install]
WantedBy=sockets.target
